# profile.d/mail
if [ -d "$HOME/Maildir" ]; then
  export MAIL="$HOME/Maildir/"
fi
